# adb-reboot/rootfs/etc/services.d/adb-reboot/finish
#!/command/with-contenv bashio

echo "adb-reboot stopped."
exit 0
