# adb-reboot/rootfs/etc/services.d/adb-reboot/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: adb-reboot
# ==============================================================================

bashio::log.info "Starting adb-reboot..."

# Make sure we use the Python interpreter from /opt/venv
exec /opt/venv/bin/python /opt/run.py
