# adb-reboot/rootfs/etc/services.d/adb-server/finish
#!/command/with-contenv bashio

echo "Stopping ADB server..."
exit 0
