# adb-reboot/rootfs/etc/services.d/adb-server/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the adb server
# ==============================================================================

bashio::log.info "Starting the ADB Server"

exec /usr/bin/adb start-server
