# adb-reboot/rootfs/etc/services.d/adb-server/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the adb server
# ==============================================================================

bashio::log.info "Starting the ADB Server"

ip a

exec /usr/bin/adb -a -P 5037 server nodaemon
